# conf/broker.ini
[program:broker]
command = /opt/rocketmq/bin/mqbroker -c  /opt/rocketmq/conf/broker.conf
stdout_logfile=/var/log/mbroker.out
