# conf/name.ini
[program:namesrv]
command = /opt/rocketmq/bin/mqnamesrv
stdout_logfile=/var/log/namesrv.out
