# conf/web.ini
[program:web]
command = java -jar /opt/rocketmq/rocketmq-console-ng-2.0.0.jar
stdout_logfile=/var/log/web.out
